# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; ------------------------------------------------------------------
; Django – migrate + collectstatic al primo avvio, poi gunicorn
; ------------------------------------------------------------------
[program:api]
directory=/api
command=/api/entrypoint.sh
environment=
    VIRTUAL_ENV="/api/.venv",
    PATH="/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ------------------------------------------------------------------
; Next.js standalone (node server.js)
; ------------------------------------------------------------------
[program:client]
directory=/client
command=node server.js
user=nextjs
environment=
    NODE_ENV="production",
    PORT="3000",
    HOSTNAME="0.0.0.0"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; ------------------------------------------------------------------
; Nginx – parte per ultimo, quando api e client sono pronti
; ------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
